# systemd/claude-session-sync.timer
[Unit]
Description=Claude Session Sync Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
